# addons/ekco/0.15.0/upgrade/ekco-upgrade.timer
[Unit]
Description=Schedule kURL auto-upgrades

[Timer]
OnCalendar=$EKCO_AUTO_UPGRADE_SCHEDULE

[Install]
WantedBy=timers.target
